# btsave-notifier-updated.service
[Unit]
Description=BTSAVE BTC Price Notifier v2
Documentation=man:node(1)
After=network.target

[Service]
Type=simple
User=xou
WorkingDirectory=/home/xou/hedge-dashboard
Environment=NODE_ENV=production
ExecStart=/usr/bin/node notifier.js
Restart=always
RestartSec=30
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=btsave-notifier

# Security (relaxed for Puppeteer)
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target